# magento-1.x/centos7/etc/supervisor/magecron.ini
[program:magecron]
command=/usr/local/bin/magecron
redirect_stderr=true
user=apache